# systemd/karipuza-bot.service
[Unit]
Description=Karipuza VPN Telegram Bot
After=network.target

[Service]
User=root
WorkingDirectory=/opt/karipuza-bot
ExecStart=/opt/karipuza-bot/venv/bin/python3 /opt/karipuza-bot/bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
